# etc/ansible/exercises/hello_world/app.conf
[program:hello_world]
command=go run /opt/site.go
autostart=false
autorestart=true
stderr_logfile=/var/log/site.err.log
stdout_logfile=/var/log/site.out.log
user=root
